# SourceCode/Java/holer-client/src/main/resources/bin/holer.service
[Unit]
Description=Holer
After=network.target

[Service]
Type=forking
ExecStart=@HOLER_HOME@/bin/startup.sh
ExecStop=@HOLER_HOME@/bin/shutdown.sh

[Install]
WantedBy=multi-user.target
